# deploy/app.service
[Unit]
Description=Zhusha XHS Studio (主程序：Node 主站 + 计费 + 抓取)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=app
WorkingDirectory=/opt/zhusha
# 主站读取 /opt/zhusha/.env 拿密钥/配置（server.js 自带 .env 读取）。
# HOME + PATH 必须带上 app 用户的 .local/bin，这样 server.js 调用本机 xhs CLI 能找到。
Environment=HOME=/home/app
Environment=PATH=/home/app/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=NODE_ENV=production
Environment=PORT=8787
ExecStart=/usr/bin/node --disable-warning=ExperimentalWarning server.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
